# systemd/pi-iot-gateway.service
[Unit]
Description=Pi IoT Gateway Service
After=network-online.target emqx.service
Wants=network-online.target
Requires=emqx.service

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pi-iot-gateway
Environment=PI_IOT_FERNET_KEY_FILE=/etc/pi-iot-gateway/fernet.key
ExecStart=/opt/pi-iot-gateway/.venv/bin/python -m backend.main
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
